# templates/uwsgi-upstart.conf.j2
# {{ ansible_managed }}

description "{{ uwsgi_upstart_description | default('') }}"
start on runlevel [2345]
stop on runlevel [06]

respawn

exec {{ _uwsgi_bin | quote }} {{ uwsgi_opts }}
